# debian/eucalyptus-cc.prerm
#!/bin/sh

if [ "$1" = "remove" ]
then
  if which invoke-rc.d >/dev/null 2>&1; then
        invoke-rc.d eucalyptus-cc stop
  else
        /etc/init.d/eucalyptus-cc stop
  fi

fi

#DEBHELPER#
