# systemd/units/eucanetd.service
[Unit]
Description=Eucalyptus network daemon
After=midolman.service network.target

[Service]
Type=forking
PIDFile=/run/eucalyptus/eucanetd.pid
ExecStart=/usr/sbin/eucanetd
LimitNOFILE=10000
LimitNPROC=100000

[Install]
WantedBy=multi-user.target

# vim: set syntax=dosini:
